# templates/heka.systemd.service.erb
[Unit]
Description=Heka event/metric/log collection and routing daemon
After=network.target auditd.service
ConditionPathExists=!/etc/heka/hekad_not_to_be_run

[Service]
EnvironmentFile=-/etc/default/heka
Type=simple
# HUP not implemented in heka yet
#ExecReload=/bin/kill -HUP $MAINPID
# http://unix.stackexchange.com/questions/263369
<% if scope['heka::forward_journald_on'] then -%>
ExecStart=-/bin/sh -c '/usr/bin/hekad -config=/etc/heka 2>&1 | /usr/bin/forward-journald -tag heka'
StandardOutput=null
StandardError=null
<% else -%>
ExecStart=-/usr/bin/hekad -config=/etc/heka
<% end -%>
ExecStartPost=/bin/sleep 3
ExecStartPost=/bin/bash -c "[ ! -z $MAINPID ]"

KillMode=process
Restart=on-failure
StandardError=inherit

[Install]
WantedBy=multi-user.target
Alias=heka.service
